# target-contrib/edge-installer/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noreset --noclear --noissue --skip-login --nohostname --autologin player - ${TERM}
StandardOutput=journal
StandardError=journal
