# target-contrib/edge-next/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noreset --noclear --noissue --skip-login --autologin player - ${TERM}
